# backend/deploy/lumilabs-backend.service
[Unit]
Description=LumiLabs unified API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
Group=www-data
WorkingDirectory=/var/www/lumilabs-backend
Environment=NODE_ENV=production
EnvironmentFile=/var/www/lumilabs-backend/.env
ExecStart=/usr/bin/env HOST=127.0.0.1 PORT=3100 /opt/lumilabs-messaging/current/bin/node server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/lumilabs-backend/uploads

[Install]
WantedBy=multi-user.target
